# launchd/com.phron.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.phron.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/terry/.cargo/bin/comes-bot</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_COMES_BOT_TOKEN</key>
        <string>{{REPLACE_WITH_TOKEN_FROM_1PASSWORD}}</string>
        <key>PHRON_SCRIPTS_DIR</key>
        <string>/Users/terry/code/phron/scripts</string>
        <key>OPENROUTER_API_KEY</key>
        <string>{{INJECT_VIA_1PASSWORD}}</string>
        <key>OPENAI_API_KEY</key>
        <string>{{INJECT_VIA_1PASSWORD}}</string>
        <key>OURA_TOKEN</key>
        <string>{{INJECT_VIA_1PASSWORD}}</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/terry/logs/comes-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/terry/logs/comes-bot.log</string>
</dict>
</plist>
